# deploy/know-canvas-caddy.service
# Know Canvas 独立 Caddy 实例 — 监听 :8081, 跟系统 caddy 互不打架
#
# 安装:
#   sudo mkdir -p /opt/know-canvas/caddy-data
#   sudo chown -R caddy:caddy /opt/know-canvas/caddy-data
#   sudo cp deploy/know-canvas-caddy.Caddyfile /opt/know-canvas/Caddyfile
#   sudo cp deploy/know-canvas-caddy.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now know-canvas-caddy
#
# 看日志:
#   sudo journalctl -u know-canvas-caddy -f

[Unit]
Description=Know Canvas - Independent Caddy Instance (:8081)
Documentation=https://caddyserver.com/docs/
After=network.target network-online.target
Requires=network-online.target

[Service]
Type=notify
User=caddy
Group=caddy
# 独立 XDG home, 避开系统 caddy 的 /var/lib/caddy storage 锁竞争
Environment=XDG_DATA_HOME=/opt/know-canvas/caddy-data
Environment=XDG_CONFIG_HOME=/opt/know-canvas/caddy-data/config
ExecStart=/usr/bin/caddy run --config /opt/know-canvas/Caddyfile --adapter caddyfile
ExecReload=/usr/bin/caddy reload --config /opt/know-canvas/Caddyfile --adapter caddyfile --address localhost:2020
TimeoutStopSec=5s
LimitNOFILE=1048576
PrivateDevices=yes
ProtectHome=true
ProtectSystem=full
AmbientCapabilities=CAP_NET_BIND_SERVICE
ReadWritePaths=/opt/know-canvas/caddy-data
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
